# install/norns/files/norns-jack.service
[Unit]
After=sound.target

[Service]
Type=simple
User=we
Group=we
Environment="JACK_NO_AUDIO_RESERVATION=1"
LimitRTPRIO=95
LimitMEMLOCK=infinity
ExecStart=/usr/bin/jackd -R -P 95 -d alsa -d hw:sndrpiproto -r 48000 -n 3 -p 128 -S -s
ExecStartPost=/usr/bin/jack_wait -w -t 10
TimeoutStopSec=1

[Install]
WantedBy=norns.target
